# wsl/supervisor/supervisor.d/i3.ini
[program:i3]
command=su - yaxin -c '/usr/bin/i3'
autostart=false
