# wsl/supervisor/supervisor.d/lightdm.ini
[program:lightdm]
command=/usr/sbin/lightdm
